# dev-tasks.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>Ngoc Tu Nguyen</Author>
    <Description>Launch dev apps at startup: VS Code, Visual Studio (Admin), Chrome, MS Teams, Spark Mail</Description>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
      <Delay>PT30S</Delay>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>cmd.exe</Command>
      <Arguments>/c start "" "C:\Users\%USERNAME%\AppData\Local\Programs\Microsoft VS Code\Code.exe"</Arguments>
    </Exec>
    <Exec>
      <Command>cmd.exe</Command>
      <Arguments>/c start "" "C:\Program Files\Microsoft Visual Studio\2022\Community\Common7\IDE\devenv.exe"</Arguments>
    </Exec>
    <Exec>
      <Command>cmd.exe</Command>
      <Arguments>/c start "" "C:\Program Files\Google\Chrome\Application\chrome.exe"</Arguments>
    </Exec>
    <Exec>
      <Command>cmd.exe</Command>
      <Arguments>/c start "" "%LOCALAPPDATA%\Microsoft\Teams\Update.exe" --processStart "Teams.exe"</Arguments>
    </Exec>
    <Exec>
      <Command>cmd.exe</Command>
      <Arguments>/c start "" "C:\Users\%USERNAME%\AppData\Local\spark\Spark.exe"</Arguments>
    </Exec>
  </Actions>
</Task>
